# etc/cron.d/aidockermon
PATH=/usr/local/bin:/usr/sbin:/usr/sbin:/usr/bin:/sbin:/bin

# Query sysload every 5s
* * * * * www-data aidockermon query sysload -r 12/5

# Query gpu every 5s
* * * * * www-data aidockermon query gpu -r 12/5

# Query disk every 1min
* * * * * www-data aidockermon query disk -r 12/5 -f / /disk

# Query containers every 1min
* * * * * dian aidockermon query containers -f DianAI

# Query sysinfo every 1day
8 8 * * * www-data aidockermon query sysinfo
