# etc/cron.d/aidockermon
PATH=/usr/local/bin:/usr/sbin:/usr/sbin:/usr/bin:/sbin:/bin

# Query sysload every 20s
* * * * * www-data aidockermon query sysload -r 3/20

# Query gpu every 20s
* * * * * www-data aidockermon query gpu -r 3/20

# Query disk every 1min
* * * * * www-data aidockermon query disk -f / /disk

# Query containers every 1min
* * * * * dian aidockermon query containers -f DianAI

# Query sysinfo every 1day
8 8 * * * www-data aidockermon query sysinfo
