# etc/cron.d/aidockermon
PATH=/home/linux/venv/aidockermon/bin:/usr/sbin:/usr/sbin:/usr/bin:/sbin:/bin

# Query sysload every 5s
* * * * * wqy aidockermon query sysload -r 12/5

# Query gpu every 5s
* * * * * wqy aidockermon query gpu -r 12/5

# Query disk every 1min
* * * * * wqy aidockermon query disk -r 12/5

# Query containers every 1min
* * * * * wqy aidockermon query containers

# Query sysinfo every 1day
8 8 * * * wqy aidockermon query sysinfo -r 12/5
